# prometheus_node_exporter/rootfs/etc/services.d/node_exporter/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Prometheus Node Exporter
# Runs the Prometheus Node Exporter
# ==============================================================================
bashio::log.info "Starting prometheus node exporter..."

# Run Prometheus
exec s6-setuidgid prometheus /usr/local/bin/node_exporter --web.config=/etc/prometheus_node_exporter_web.yml